# deploy/supervisor/polariss.conf
[program:polariss]
command= /home/jair/polariss/app/deploy/gunicorn.sh
user=jair
stdout_logfile= /home/jair/polariss/app/logs/supervisor.log
redirect_stderr= true
autostart= true
autorestart= true
environment=LANG= en_US.UTF-8,LC_ALL=en_US.UTF-8
